# etc/profile.d/arango.alias.sh


alias csnbopsh='docker run --rm -it arangodb:3.4 arangosh --server.endpoint "tcp://csnbop:8529" --server.username root --server.database csnhub'

alias csnboesh='docker run --rm -it arangodb:3.4 arangosh --server.endpoint "tcp://csnboe:8529" --server.username root --server.database snmtest'

