# etc/profile.d/borg.alias.sh
## Borg Backup 
## gcugli Dic2018

alias borgc='borg create -s -p -C auto,lzma'
alias borgi='borg init -e repokey'

