# etc/profile.d/cgc.alias.sh
## gcugli Sep2017

alias qrclip="xclip -selection clip‐board -out | qrencode -s10 -o - | display -"

#### System
alias nore='cut -b 1-$COLUMNS | more'




