# etc/profile.d/docker.alias.sh
## Docker alias
## gcugli Sep2017


alias dk='docker'
alias dkc='docker-compose'
alias dkn='docker network'
alias dkv='docker volume'
alias dki='docker inspect'
alias dkps='docker ps --format "table {{.ID}}  {{.Names}}\t{{.Image}}\t{{.Status}}"'


