# etc/profile.d/docker.sh
# docker aliases and functions

export RPROXYDOMAIN=botnet.ar

alias dk=docker
alias dkc=docker-compose

