# etc/profile.d/rsync.alias.sh


alias rsyncar='rsync -avzh --stats --progress'


