# etc/profile.d/syctl.alias.sh
## Systemd
## gcugli Sep2017

alias syctl='sudo systemctl'
alias syctl.s='sudo systemctl status'
alias syctl.d='sudo systemctl stop'
alias syctl.u='sudo systemctl start'
alias syctl.r='sudo systemctl daemon-reload'
alias sylog='sudo journalctl -u'
alias sylogf='sudo journalctl -fu'


